# src-tauri/notifierd/rencal-notifierd.service
[Unit]
Description=Rencal calendar reminder notifier
Documentation=https://github.com/anthropics/rencal/blob/main/docs/notifications.md
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/rencal-notifierd
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=graphical-session.target
